# festy.service.example
[Unit]
Description=festy Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace <user> with your Linux username and adjust paths to where festy lives.
User=<user>
WorkingDirectory=/home/<user>/festy
ExecStart=/home/<user>/festy/.venv/bin/python /home/<user>/festy/bot.py
EnvironmentFile=/home/<user>/festy/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target